# examples/chatbot/supervisor/agentforge.conf
[unix_http_server]
file = /tmp/supervisor.sock

[supervisord]
logfile = /tmp/supervisord.log
pidfile = /tmp/supervisord.pid
childlogdir = /tmp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:agentforge-chatbot]
command = cd /home/ubuntu/workspace/agentforge/examples/chatbot/backend && uvicorn main:app --host 0.0.0.0 --port 8000
directory = /home/ubuntu/workspace/agentforge/examples/chatbot/backend
user = ubuntu
autostart = true
autorestart = true
stdout_logfile = /tmp/agentforge-chatbot-stdout.log
stderr_logfile = /tmp/agentforge-chatbot-stderr.log
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
environment = PYTHONPATH="/home/ubuntu/workspace/agentforge/src:%(ENV_PYTHONPATH)s"