# projects/chatbot/supervisor/chatbot.conf
[program:agentforge-chatbot]
command = /home/ubuntu/workspace/agentforge/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
directory = /home/ubuntu/workspace/agentforge/projects/chatbot/backend
user = ubuntu
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /tmp/agentforge-chatbot-stdout.log
stdout_logfile_maxbytes = 10MB
environment = PYTHONPATH="/home/ubuntu/workspace/agentforge/src"